# systemd/claude-usage.timer
[Unit]
Description=Claude Code usage refresh (every 1 min)
Requires=claude-usage.service

[Timer]
OnBootSec=30
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=false

[Install]
WantedBy=timers.target
