# systemd/codex-usage.timer
[Unit]
Description=OpenAI Codex usage refresh (every 5 min)
Requires=codex-usage.service

[Timer]
OnBootSec=45
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=false

[Install]
WantedBy=timers.target
