# systemd/opencode-usage.timer
[Unit]
Description=OpenCode usage refresh (every 1 min)
Requires=opencode-usage.service

[Timer]
OnBootSec=45
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=false

[Install]
WantedBy=timers.target
