# systemd/opencode-usage.timer
[Unit]
Description=OpenCode usage refresh (every 5 min)
Requires=opencode-usage.service

[Timer]
OnBootSec=45
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=false

[Install]
WantedBy=timers.target
